# systemd/jasper-webui.service
[Unit]
Description=Jasper Speaker Web Dashboard
After=squeezelite.service network-online.target
Wants=network-online.target

[Service]
User=jaspertech
Group=jaspertech
WorkingDirectory=/home/jaspertech/jasper-speaker/webui
ExecStart=/home/jaspertech/jasper-speaker/webui/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
