# systemd/jasper-docker.service
[Unit]
Description=Jasper Docker Services (Music Assistant + Home Assistant)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=jaspertech
WorkingDirectory=/home/jaspertech/jasper-speaker/docker
ExecStart=/usr/bin/docker compose up -d
ExecStop=/usr/bin/docker compose down
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
